# templates/sqs2slack.ini
[program:sqs2slack]
user={{ sqs2slack.user }}
environment=PATH="PATH:/usr/local/bin:/usr/bin:/bin",HOME="/home/{{ sqs2slack.user }}"
autorestart = True
log_stderr = True
stopsignal = QUIT
log_stdout = True
command = sh /var/deploy/sqs2slack.sh
autostart = True
stderr_logfile = /var/log/supervisor/sqs2slack.err.log
stdout_logfile = /var/log/supervisor/sqs2slack.stdout.log

